# deploy/crypto_promo_bot.service
[Unit]
Description=Crypto Promo Bot - Telegram Bot for Crypto Promotions
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/crypto_promo_bot
Environment=PATH=/opt/crypto_promo_bot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/crypto_promo_bot/venv/bin/python main.py

# Автоматический перезапуск при сбоях
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto_promo_bot

[Install]
WantedBy=multi-user.target
